# configs/docker_dashboard.conf
[program:docker_dashboard]
directory=/usr/lib/docker-dashboard
command=/usr/lib/docker-dashboard/env/bin/uvicorn server.asgi:application --port 2121

autostart=true
autorestart=true

stdout_logfile=/usr/lib/docker-dashboard/docker_dashboard.log
stderr_logfile=/usr/lib/docker-dashboard/docker_dashboard.log.error


[inet_http_server]
port=127.0.0.1:9001